# curator4-2029619917.conf
[program:curator4-2029619917]
command=python3 /home/curator4/curator4-2029619917/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-2029619917/test.err.log
stdout_logfile=/home/curator4/curator4-2029619917/test.out.log